# planning-poker-bot.service
[Unit]
Description=Planning Poker Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=bot
Group=bot
WorkingDirectory=/opt/planning-poker-bot
ExecStart=/usr/bin/python3 bot.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/planning-poker-bot
Environment=PYTHONUNBUFFERED=1

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/planning-poker-bot/data

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
