# telegram-pb-bot.service
[Unit]
Description=Telegram Planning Poker Bot
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/telegram_pb
# Load environment from .env file (create it from env.example)
EnvironmentFile=/path/to/telegram_pb/.env
ExecStart=/usr/bin/python3 /path/to/telegram_pb/run.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
