# telegram-pb-bot.service
[Unit]
Description=Telegram Planning Poker Bot
After=network.target

[Service]
Type=simple
User=paulshults
WorkingDirectory=/Users/paulshults/dev/telegram_pb
Environment=BOT_TOKEN=7826785109:AAGIgiAvvrFVR7kYKuQa6z4WqhYzwGjvGhg
Environment=USER_TOKEN=user_token
Environment=LEAD_TOKEN=lead_token
Environment=ADMIN_TOKEN=admin_token
Environment=JIRA_URL=https://media-life.atlassian.net
Environment=JIRA_USERNAME=pe_shults@betboom.org
Environment=JIRA_API_TOKEN=ATATT3xFfGF0wzdejcmqllyiWtn21dLG7mKsBQ02kfZrbICt4ebrCmdSAAsxewmHpwvVwhIP7mYxJipPGz_ACa3c_9_OlvcNnw1mF6erKszwVHnXMQGVLK6u2xW36OVAEVobB2oqeIHSt1LfTocmR4ewYI37_zkLG3KnQ-e9NAamwPtc1Z8MWjQ=E13ED861
Environment=STORY_POINTS_FIELD=customfield_10022
Environment=SUPPORTED_TOPICS={"-1003087077812": ["ALL"]}
Environment=STATE_FILE=data/state.json
ExecStart=/usr/bin/python3 /Users/paulshults/dev/telegram_pb/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
